# syncp/syncp.service
# syncp.service — KeePass vault backup pipeline (systemd user unit).
#
# Layout this unit assumes:
#   ~/.local/bin/scripts/syncp.py    — the pipeline script (stowed from
#                                      dotfiles/scripts/.local/bin/scripts/)
#   ~/.secrets/syncp.env              — KEY=VALUE file with BORG_PASSPHRASE=... RCLONE_CONFIG_PASS=...
#                                      (mode 0600, dir 0700)
#   ~/systemdcustom/syncp/*.service  — this repo
#   ~/systemdcustom/syncp/*.timer
#
# Install:
#   1. Make sure the secret file exists and is locked down:
#        install -m 700 -d ~/.secrets
#        printf 'BORG_PASSPHRASE=%s\n' 'your-passphrase' > ~/.secrets/syncp.env
#        printf 'RCLONE_CONFIG_PASS=%s\n' 'your-passphrase' > ~/.secrets/syncp.env
#        chmod 600 ~/.secrets/syncp.env
#   2. Symlink (or copy) the unit files into the user units directory.
#        Symlinking lets you `git pull` upstream changes without re-copying.
#        mkdir -p ~/.config/systemd/user
#        ln -s ~/systemdcustom/syncp/syncp.service        ~/.config/systemd/user/
#        ln -s ~/systemdcustom/syncp/syncp.timer          ~/.config/systemd/user/
#        ln -s ~/systemdcustom/syncp/syncp-failure.service ~/.config/systemd/user/
#   3. Enable the timer (it pulls in the service when it fires):
#        systemctl --user daemon-reload
#        systemctl --user enable --now syncp.timer
#   4. Test once by hand:
#        systemctl --user start syncp.service
#        journalctl --user -u syncp.service -e
#
# Pipeline details and exit codes are in the script's module docstring.

[Unit]
Description=KeePass vault → Borg → cloud-remote backup
Documentation=file:%h/.local/bin/scripts/syncp.py

# Fire syncp-failure.service if we exit non-zero (i.e. on EXIT_LOCAL_FAILURE
# or EXIT_REMOTE_FAILURE). The "skipped because offline" path returns 0,
# so a flaky network on a given day will NOT trigger the notifier.
OnFailure=syncp-failure.service

[Service]
Type=oneshot

# Block until NetworkManager has a routable connection, up to 120s. In a
# user unit, After=network-online.target is a no-op (that target lives in
# the system manager), so without this probe the service can fire seconds
# after boot/login — before Wi-Fi has associated — and the script's offline
# path silently skips the cloud sync. nm-online exits 0 once online and
# non-zero on timeout; on timeout we still ExecStart the script, which
# will hit its own offline path and exit 0 cleanly.
# ExecStartPre=-/usr/bin/nm-online -q --timeout=120
ExecStart=%h/.local/bin/scripts/syncp.py

# Hard ceiling on the entire run. Must be larger than the script's
# internal timeouts (RCLONE_TIMEOUT × len(REMOTES) ≈ 6h) so syncp.py's
# own per-command timeout fires first and we get a clean error log
# rather than systemd killing the process mid-rclone.
TimeoutStartSec=8h

# Source the Borg passphrase from a KEY=VALUE file. The file lives
# outside the public dotfiles repo and is mode 0600. systemd reads it
# at service start and injects every entry into the unit's environment,
# so the script sees BORG_PASSPHRASE the same as if it had been exported
# in the parent shell.
#
# Trade-off vs. LoadCredential=: the passphrase is now visible in
# /proc/<pid>/environ to root and to processes running as this user.
# That's fine for a single-user laptop; on a shared box, prefer
# LoadCredential= + BORG_PASSCOMMAND / RCLONE_PASSWORD_COMMAND.
# Inject the Homebrew bin path so systemd can find python3, it's
# libraries, rclone and borg
Environment="PATH=/home/linuxbrew/.linuxbrew/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=%h/.secrets/syncp.env

# Be a polite background job: never preempt interactive work.
Nice=10
IOSchedulingClass=idle

# Sandboxing. Kept modest because this is a user unit that legitimately
# needs to read ~/vault and write ~/keepass-backups — ProtectHome would
# break both. The kernel-side restrictions below are free wins.
NoNewPrivileges=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Note: no [Install] section. The service is pulled in by syncp.timer;
# enabling it directly would bypass the schedule.
